# lab/provision/osquery.conf
{
  "options": {
    "logger_plugin": "filesystem",
    "logger_path": "/var/log/osquery",
    "schedule_splay_percent": "10",
    "utc": "true"
  },
  "schedule": {
    "crontab_snapshot": {
      "query": "SELECT command, path, event FROM crontab;",
      "interval": 300,
      "snapshot": true,
      "description": "Cron persistence visibility (T1053.003)"
    },
    "listening_ports": {
      "query": "SELECT p.name, p.path, lp.port, lp.address, lp.protocol FROM listening_ports lp JOIN processes p USING (pid) WHERE lp.port != 0;",
      "interval": 300,
      "description": "Differential view of listening services"
    },
    "users_snapshot": {
      "query": "SELECT username, uid, gid, shell, directory FROM users;",
      "interval": 600,
      "snapshot": true,
      "description": "New/anomalous user detection input (T1136)"
    },
    "systemd_units": {
      "query": "SELECT id, description, fragment_path, active_state FROM systemd_units WHERE fragment_path LIKE '/etc/systemd/system/%';",
      "interval": 600,
      "description": "Persistence via systemd service (T1543.002)"
    }
  }
}
